# launchd/com.desmond.codexremote.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.desmond.codexremote</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__ROOT_DIR__/scripts/start_server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ROOT_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Applications/Codex.app/Contents/Resources:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.codexremote/service.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.codexremote/service.err.log</string>
  </dict>
</plist>
